# .config/systemd/user/tmp.service
[Unit]
Description=Remove tmp dir

[Service]
ExecStart=/usr/bin/bash -c "/usr/bin/rm -rfv %h/tmp/*"

[Install]
WantedBy=mystuff.target
