# tag-xps13/config/systemd/user/wallpaper.timer
[Unit]
Description=Periodically change wallpaper

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h
AccuracySec=50s

[Install]
WantedBy=mystuff.target
